# templates/userli-collector.timer.j2
[Unit]
Description=Userli Prometheus Textfile Collector Timer

[Timer]
OnCalendar={{ userli_prometheus_textfile_collector_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
